# udev/70-persistent-display.rules
# rule to run autorandr when display is plugged in
SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/lemmsh/script/autor.sh"
